# release/systemd/elyan.service.example
[Unit]
Description=Elyan local-first personal agent
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=elyan
Group=elyan
EnvironmentFile=/srv/elyan/.env
WorkingDirectory=/srv/elyan/current
ExecStart=/usr/bin/node /srv/elyan/current/.next/standalone/server.js
Restart=always
RestartSec=5
TimeoutStopSec=30

[Install]
WantedBy=multi-user.target
